# packaging/tizenclaw-code-sandbox.service
[Unit]
Description=TizenClaw Code Execution Sandbox
After=local-fs.target
Before=tizenclaw.service

[Service]
Type=simple
ExecStart=/usr/libexec/tizenclaw/tizenclaw_secure_container.sh start
ExecStop=/usr/libexec/tizenclaw/tizenclaw_secure_container.sh stop
Restart=on-failure
RestartSec=3
User=root

[Install]
WantedBy=multi-user.target
